# files/shared/etc/profile.d/nova-aliases.sh
#!/usr/bin/env bash

alias fastfetch="/usr/bin/fastfetch --config /usr/share/nova/configs/config.jsonc"
